# debian/prerm
#!/bin/sh
set -e

case "$1" in
  remove|upgrade|deconfigure)
    systemctl stop otzi otzi-relay 2>/dev/null || true
    systemctl disable otzi otzi-relay 2>/dev/null || true
    ;;
esac

#DEBHELPER#
